# services/99-umbrel-backup.rules
# Umbrel Guardian — auto-mount backup drive when plugged in.
# Triggers the mount script whenever a partition labeled "umbrel-backup" appears.
ACTION=="add", ENV{ID_FS_LABEL}=="umbrel-backup", RUN+="/usr/bin/systemd-run --no-block /usr/local/bin/mount-umbrel-backup.sh"
