# services/umbrel-guardian-bot.service
[Unit]
Description=Umbrel Guardian Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=umbrel
WorkingDirectory=/home/umbrel/umbrel/umbrel-guardian
# Brief delay so DNS/networking is fully ready after network-online.target
ExecStartPre=/bin/sleep 5
ExecStart=/home/umbrel/umbrel/umbrel-guardian/.venv/bin/python3 /home/umbrel/umbrel/umbrel-guardian/bot/guardian_bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Hardening
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/home/umbrel/umbrel/umbrel-guardian
MemoryMax=128M
CPUQuota=20%
# NOTE: NoNewPrivileges is intentionally NOT set — the bot uses `sudo -n` to
# invoke scripts/system_control.sh for /system_reboot, /system_shutdown,
# /restart_docker, /restart_umbrel. Sudo needs setuid to escalate to root,
# which NoNewPrivileges=yes blocks. The privilege boundary is enforced by
# /etc/sudoers.d/umbrel-guardian-system, which permits only the five exact
# subcommands of system_control.sh — nothing else.

[Install]
WantedBy=multi-user.target
