# services/umbrel-guardian-daily.timer
[Unit]
Description=Umbrel Guardian Daily Status Summary Timer

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
