# services/umbrel-guardian-health.timer
[Unit]
Description=Umbrel Guardian Health Check Timer

[Timer]
OnCalendar=*:0/30
Persistent=true

[Install]
WantedBy=timers.target
